# debian/setup-modem.init
#!/sbin/openrc-run

command="/usr/sbin/setup-modem.sh"
command_background="true"

depend() {
	after eg25
	before ofono
}
